# tests/docker/configs/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:rbee-hive]
command=/home/rbee/.local/bin/rbee-hive --port 9000
user=rbee
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/home/rbee",USER="rbee",RBEE_CONFIG_DIR="/home/rbee/.config/rbee",RBEE_DATA_DIR="/home/rbee/.local/share/rbee"
